# contrib/systemd/satd.service
# satd.service — systemd unit for the satd Bitcoin node.
#
# Install:
#   sudo install -Dm644 contrib/systemd/satd.service /etc/systemd/system/satd.service
#   sudo install -Dm755 target/release/satd /usr/local/bin/satd
#   sudo install -Dm755 target/release/sat-cli /usr/local/bin/sat-cli
#   sudo useradd --system --home /var/lib/satd --shell /usr/sbin/nologin satd
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now satd
#
# Cookie-auth access for non-root operators: add yourself to the
# `satd` group. The unit chmods $DATADIR/.cookie to 0640 on every
# start, so group members can run `sat-cli` / `sat-tui` without sudo:
#   sudo usermod -aG satd <your-user>
# (group changes only take effect in new login sessions)
#
# Optional config (recommended): drop a `bitcoin.conf` (or `satd.conf`)
# into `$SATD_DATADIR/bitcoin.conf` (default `/var/lib/satd/bitcoin.conf`).
# satd reads it via the standard `-conf=` resolution against the datadir.
#
# Runtime overrides — `/etc/default/satd` (optional):
#   SATD_DATADIR=/var/lib/satd       # override datadir path
#   SATD_OPTS="--rpcport=18880 ..."  # extra flags appended to ExecStart
#
# If you set SATD_DATADIR to anything other than `/var/lib/satd`, you
# also need a drop-in that updates the hardening directives, since
# systemd does not expand environment variables in those:
#   /etc/systemd/system/satd.service.d/datadir.conf:
#     [Service]
#     StateDirectory=
#     ReadWritePaths=
#     ReadWritePaths=/your/datadir
# (Empty values reset the inherited list; the second line adds yours.)
#
# Type=notify: satd calls sd_notify(READY=1) only after every listener
# is bound (RPC, P2P, optional Esplora/Electrum/MCP/events surfaces),
# so dependent units (Tor onion services, watchtower processes) start
# at the right moment instead of racing satd's bind sequence.
#
# During reindex (which can run for many hours on a fully-synced
# mainnet node) satd emits sd_notify(EXTEND_TIMEOUT_USEC=120000000)
# every 30s plus STATUS=<phase: progress>. TimeoutStartSec is `infinity`
# because the heartbeat IS the liveness check — silence for >120s and
# systemd kills the unit (correctly: it's stuck). `systemctl status
# satd` shows the live STATUS line during reindex.

[Unit]
Description=satd — Bitcoin full node
Documentation=https://github.com/epochbtc/satd/blob/master/docs/PACKAGING.md
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main

# Operators can override the datadir or append extra flags via
# /etc/default/satd without forking the unit. ${SATD_DATADIR} stays a
# single argv token; $SATD_OPTS is word-split, so multi-flag strings
# like `--rpcport=18880 --txindex` expand to separate arguments.
Environment=SATD_DATADIR=/var/lib/satd
Environment=SATD_OPTS=
EnvironmentFile=-/etc/default/satd
ExecStart=/usr/local/bin/satd --datadir=${SATD_DATADIR} $SATD_OPTS
ExecStop=/bin/kill -SIGTERM $MAINPID

# Loosen the cookie file from 0600 → 0640 once it's been written, so
# users added to the `satd` group can authenticate sat-cli/sat-tui
# without sudo. `-` prefix: missing cookie (e.g., when --rpcuser is
# set instead) is non-fatal. Cookie still rotates on every daemon
# start, so this is no less secure than 0600 for outside-group users.
ExecStartPost=-/bin/chmod 0640 ${SATD_DATADIR}/.cookie

# Reindex can take hours. The EXTEND_TIMEOUT_USEC heartbeats inside
# satd reset this every 30s during the pre-bind phase; `infinity`
# removes the static budget that would otherwise kill us during
# `--reindex-chainstate`. See the comment block above.
TimeoutStartSec=infinity

# RocksDB flush + undo-file fsync on shutdown can take a while under load.
# 10 minutes matches Bitcoin Core's bitcoind.service ceiling and is the
# practical upper bound observed during mainnet-tip operation.
TimeoutStopSec=10min
KillSignal=SIGTERM
SendSIGKILL=yes

# Runtime liveness: systemd kills the unit if it goes 60s without a
# WATCHDOG=1 ping. After notify_ready, satd spawns a watchdog task that
# ticks at WATCHDOG_USEC/2 (= 30s here) and pings only when its
# subsystem probes report healthy (currently: a non-blocking try_read
# on the chain-state tip lock). Catches runtime-level wedges that the
# chain-level stall_watchdog (which itself runs on the tokio runtime)
# can't detect on its own. Restart=always (below) brings us back after
# a watchdog kill.
#
# 60s is the tightest interval that survives a brief scheduler stall
# (e.g. a long RocksDB compaction pause) without false-positive kills.
# Tightening this further requires per-subsystem probes proven bounded
# under load.
WatchdogSec=60s

# `Restart=always` (vs. on-failure) is deliberate: the in-process stall
# watchdog (node/src/stall_watchdog.rs) escalates a detected wedge by
# raising SIGTERM, which runs the graceful `--max-shutdown-secs` flush
# path and exits 0. With `Restart=on-failure` that watchdog-triggered
# shutdown would NOT auto-restart — issue #178 had a workstation IBD
# wedge SIGTERM cleanly and then sit `inactive (dead)` for 3+ hours.
# Operator-driven `systemctl stop` still stops the unit; only
# unsupervised exits (watchdog SIGTERM, abort path, OOM kill) get
# auto-restarted.
Restart=always
RestartSec=5s

User=satd
Group=satd

# StateDirectory creates /var/lib/satd 0750 satd:satd if it doesn't already
# exist. Survives package reinstalls; preserved across stop/start.
StateDirectory=satd
StateDirectoryMode=0750

# Resource limits.
#  - LimitNOFILE: RocksDB keeps many SST files mmapped; 65536 leaves
#    headroom over the ~thousands the chainstate uses at mainnet tip.
#  - TasksMax: no thread bomb, but rayon + tokio + RPC hyper threads
#    add up. 16k is plenty without being unbounded.
LimitNOFILE=65536
TasksMax=16384

# Hardening. These are deliberately conservative — satd needs network
# (P2P + outbound HTTPS), filesystem write to its datadir, and that's
# the lot.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
# SystemCallFilter is additive: the first line establishes the
# allowed base set (@system-service, the broad "normal daemon" set
# from systemd's predefined groups). The second line uses the `~`
# prefix to subtract groups from that set — here, @privileged
# (capability-related calls like setuid) and @resources (rlimit /
# scheduler tweaks). Net result: a daemon-shaped allowlist with the
# admin/tuning subsets carved back out.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/satd

[Install]
WantedBy=multi-user.target
